# configs/crontab
@daily /usr/local/bin/bekup /data/config.json