# configs/crontab
@daily /usr/local/bin/bekup --config=/data/config.json